# scripts/com.openwalnut.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Open Walnut watchdog LaunchAgent template.
  Install via: bash scripts/install-watchdog.sh
  The install script substitutes __WALNUT_DIR__ with the absolute path.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openwalnut.watchdog</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__WALNUT_DIR__/scripts/walnut-watchdog.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>WALNUT_DIR</key>
        <string>__WALNUT_DIR__</string>
        <key>WALNUT_PORT</key>
        <string>3456</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/tmp/open-walnut/watchdog.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/open-walnut/watchdog.stderr.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
